# ansible/roles/server_initial_setup/files/ethtool-ring-buffers.service
[Unit]
Description=Set physical NIC ring buffers to maximum
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/set-ethtool-ring-buffers.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
