# ansible/roles/solana_validator_agave/templates/validator-logrotate.service.j2
[Unit]
Description=Logrotate for validator logs

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate -f /etc/logrotate.d/{{ validator_service_name }}
