# ansible/roles/solana_validator_agave/templates/validator-logrotate.timer.j2
[Unit]
Description=Run logrotate for {{ validator_service_name }} every hour
Requires=validator-logrotate-{{ validator_service_name }}.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
