# ansible/roles/monitoring_hub/templates/influxdb/influxdb.service.j2
# Managed by Ansible - do not edit manually

[Unit]
Description=InfluxDB is an open-source, distributed, time series database
Documentation=https://docs.influxdata.com/influxdb/
After=network-online.target
Wants=network-online.target

[Service]
User=influxdb
Group=influxdb
LimitNOFILE=65536
EnvironmentFile=-/etc/default/influxdb2
ExecStart=/usr/bin/influxd $INFLUXD_OPTS
Type=simple
KillMode=control-group
Restart=on-failure
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
Alias=influxd.service
