# ansible/roles/monitoring_hub/templates/validator_metrics_service/validator-metrics.service.j2
[Unit]
Description=Solana Validator Metrics Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Load environment variables from centralized file
EnvironmentFile=-/usr/local/etc/validator-metrics.env
ExecStart=/bin/bash -c 'while true; do /usr/local/bin/validator-metrics.sh; sleep {{ solana_metrics_interval }}; done'
Restart=always
RestartSec=5
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
User=root

[Install]
WantedBy=multi-user.target